# examples/audio/vupeak.service
[Unit]
Description=vupeak service

[Service]
ExecStart=/usr/local/bin/vupeak_wrapper
Restart=always

[Install]
WantedBy=multi-user.target
