# examples/olaclient/olaclient.service
[Unit]
Description=olaclient service
After=olad.service

[Service]
User=olad
ExecStart=/bin/sh -c '/usr/local/bin/olaclient 2>&1 >> /var/log/olaclient.log'
Restart=always

[Install]
WantedBy=multi-user.target
