# examples/olaclient/olac.service
[Unit]
Description=olaclient service

[Service]
ExecStart=/usr/local/bin/olac_wrapper

[Install]
WantedBy=default.target
